# USBbackup.service
[Unit]
Description=USB media watcher for automatic backup
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/bash "%h/.local/bin/USBbackup.sh"
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

